# scripts/systemd/adom-bot.service
[Unit]
Description=ADOM Flashloan Bot (Python)
After=network.target atom-backend.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/atom
EnvironmentFile=/etc/atom/backend.env
ExecStart=/usr/bin/python3 backend/bots/working/atom_hybrid_bot.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/atom/adom-bot.log
StandardError=append:/var/log/atom/adom-bot.err.log

[Install]
WantedBy=multi-user.target

