# scripts/aeon.service
[Unit]
Description=AEON Arbitrage Backend
After=network.target

[Service]
User=root
WorkingDirectory=/root/arbitrage-trustless-onchain-module
ExecStart=/root/arbitrage-trustless-onchain-module/venv/bin/python3 /root/arbitrage-trustless-onchain-module/start_aeon_bulletproof.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
